# scripts/bcm2711/99-camera-setup.rules
# debian/99-camera-setup.rules
ACTION=="add", SUBSYSTEM=="video4linux", KERNEL=="v4l-subdev[0-9]*", RUN+="/usr/local/bin/set_rpi4_pipeline %k"
